# supervisor/require_armada_stats.conf
[program:require-armada-stats-influxdb]
directory=/opt/microservice/src/local_magellan
command=python require_service.py 32776 armada-stats-influxdb

